# 5.4.43/config/supervisord-phpfpm5.4.conf
[program:phpfpm5.4]
command=/opt/scripts/start-phpfpm-5.4.43.sh
numprocs=1
autostart=true
autorestart=true
