# 5.6.16/config/supervisord-phpfpm5.6.conf
[program:phpfpm5.6]
command=/opt/scripts/start-phpfpm-5.6.16.sh
numprocs=1
autostart=true
autorestart=true
