# 7.0.0/config/supervisord-phpfpm7.0.conf
[program:phpfpm7.0]
command=/opt/scripts/start-phpfpm-7.0.0.sh
numprocs=1
autostart=true
autorestart=true
