# 5.5.30/config/supervisord-phpfpm5.5.conf
[program:phpfpm5.5]
command=/opt/scripts/start-phpfpm-5.5.30.sh
numprocs=1
autostart=true
autorestart=true
