# 5.3.29/config/supervisord-phpfpm5.3.conf
[program:phpfpm5.3]
command=/opt/scripts/start-phpfpm-5.3.29.sh
numprocs=1
autostart=true
autorestart=true
